# IntelIntent-Seed/IntelIntent_Mutation_Launcher (1).xml
<?xml version='1.0' encoding='utf-8'?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2"><RegistrationInfo><Date>2023-01-01T12:00:00</Date><Author>IntelIntent</Author></RegistrationInfo><Triggers><BootTrigger><Enabled>true</Enabled></BootTrigger></Triggers><Principals><Principal id="Author"><RunLevel>HighestAvailable</RunLevel><LogonType>InteractiveToken</LogonType><UserId>S-1-5-18</UserId></Principal></Principals><Settings><MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy><DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries><StopIfGoingOnBatteries>false</StopIfGoingOnBatteries><AllowHardTerminate>true</AllowHardTerminate><StartWhenAvailable>true</StartWhenAvailable><RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable><Enabled>true</Enabled><Hidden>false</Hidden><ExecutionTimeLimit>PT1H</ExecutionTimeLimit><Priority>7</Priority></Settings><Actions Context="Author"><Exec><Command>powershell.exe</Command><Arguments>-ExecutionPolicy Bypass -File "C:\IntelIntent_PC_Reset\run_all_modules.ps1"</Arguments></Exec></Actions></Task>